# .setup/supervisor/supervisor.conf
[supervisord]
daemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
autostart=true
autorestart=true
startsecs=5
startretries=3
numprocs=1
numprocs_start=1
priority=10
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log
command=/usr/sbin/nginx -g 'daemon off;'

[program:php-fpm]
autostart=true
autorestart=true
startsecs=5
startretries=3
numprocs=1
numprocs_start=1
priority=5
stderr_logfile=/var/log/nginx/php-error.log
stdout_logfile=/var/log/nginx/php-access.log
command=/usr/local/sbin/php-fpm -F
